# rootfs/target/ftm-100s/etc/init.d/hostapd
#!/bin/sh /etc/rc.common
# (C) 2008 openwrt.org

START=50

start() {
	hostapd /etc/hostapd.conf &
}

stop() {
	killall hostapd
}
