# rootfs/target/ftm-100t-tp/etc/crontabs/root
0 2 * * * rdate -s 203.248.240.140 && hwclock -w
