# rootfs/target/tpgw/opt/thingplus-gateway/device/update/bone/files/wvdial.crontab
*/5 * * * * root bash /etc/init.d/wvdial.sh start
