# rootfs/base/usr/lib/opkg/info/busybox.list
/bin/dmesg
/usr/bin/id
/bin/login
/bin/rm
/usr/bin/wget
/sbin/lsmod
/bin/sh
/sbin/vconfig
/usr/bin/test
/sbin/mkdosfs
/sbin/depmod
/bin/zcat
/bin/fgrep
/bin/addgroup
/sbin/halt
/bin/ash
/usr/sbin/flashcp
/sbin/reboot
/bin/true
/sbin/route
/usr/bin/md5sum
/usr/bin/free
/sbin/logread
/usr/bin/time
/usr/bin/passwd
/sbin/klogd
/usr/bin/traceroute
/usr/sbin/rdev
/bin/more
/sbin/tune2fs
/bin/vi
/usr/bin/tftp
/sbin/modprobe
/bin/pwd
/usr/bin/nc
/sbin/mke2fs
/usr/bin/renice
/bin/chown
/bin/delgroup
/etc/init.d/cron
/usr/sbin/rdate
/bin/mount
/bin/sync
/usr/bin/dirname
/usr/bin/bzip2
/sbin/rmmod
/usr/bin/hostid
/sbin/hwclock
/sbin/udhcpc
/sbin/insmod
/bin/ping6
/usr/bin/basename
/bin/echo
/usr/bin/less
/usr/sbin/chroot
/usr/sbin/crond
/usr/bin/ftpget
/usr/bin/head
/bin/gunzip
/bin/netstat
/usr/sbin/ftpd
/bin/busybox
/bin/grep
/usr/bin/clear
/usr/bin/seq
/bin/sleep
/bin/nice
/usr/bin/diff
/sbin/watchdog
/bin/kill
/bin/mktemp
/usr/bin/tr
/usr/bin/mesg
/sbin/mkfs.minix
/bin/mknod
/usr/bin/resize
/bin/rmdir
/usr/bin/uniq
/usr/bin/tee
/usr/bin/ftpput
/usr/bin/bzcat
/usr/bin/hexdump
/usr/sbin/flash_unlock
/bin/date
/usr/bin/wc
/usr/bin/nslookup
/usr/bin/[
/usr/bin/pgrep
/bin/umount
/sbin/swapon
/bin/netmsg
/usr/bin/tail
/bin/stat
/bin/mountpoint
/usr/bin/pkill
/usr/bin/mkfifo
/bin/chgrp
/sbin/mdev
/sbin/swapoff
/bin/cp
/sbin/poweroff
/usr/bin/[[
/bin/dd
/bin/lock
/bin/df
/bin/false
/usr/bin/killall
/usr/bin/sort
/bin/mkdir
/sbin/ifconfig
/usr/bin/cmp
/usr/bin/reset
/sbin/start-stop-daemon
/usr/sbin/dhcprelay
/bin/lsattr
/bin/usleep
/bin/cat
/usr/bin/length
/usr/sbin/flash_eraseall
/sbin/mkfs.vfat
/etc/init.d/telnet
/usr/bin/arping
/bin/deluser
/usr/bin/telnet
/sbin/getty
/bin/touch
/sbin/mkswap
/bin/tar
/usr/bin/taskset
/usr/bin/lsusb
/sbin/fsck
/usr/bin/yes
/usr/bin/top
/sbin/fdisk
/sbin/switch_root
/usr/bin/expr
/usr/bin/env
/sbin/mkfs.ext2
/sbin/sysctl
/sbin/syslogd
/bin/sed
/usr/bin/logger
/bin/chmod
/usr/bin/awk
/bin/uname
/sbin/blkid
/usr/bin/which
/bin/chattr
/usr/bin/find
/bin/ln
/bin/ls
/sbin/pivot_root
/usr/bin/nohup
/usr/bin/cut
/usr/bin/du
/bin/mv
/usr/sbin/telnetd
/usr/sbin/flash_lock
/usr/bin/crontab
/usr/bin/printf
/usr/bin/lspci
/sbin/init
/usr/bin/tty
/bin/ping
/sbin/hdparm
/usr/bin/bunzip2
/sbin/mkfs.reiser
/usr/bin/uptime
/usr/bin/killall5
/bin/fsync
/bin/pidof
/usr/bin/strings
/bin/egrep
/usr/sbin/udhcpd
/sbin/arp
/usr/bin/patch
/bin/ps
/bin/adduser
/bin/cpio
/usr/bin/xargs
/usr/bin/hd
/bin/gzip
/usr/bin/dumpleases
