# rootfs/base/usr/lib/opkg/info/kmod-atheros-11AC-10.1.357-wfo.list
/lib/modules/2.6.36/atheros_wfo/umac.ko
/sbin/killVAP
/etc/ath/WSC_sta.conf
/lib/modules/2.6.36/wireless
/etc/ath/WPA-FT-PSK.ap_bss
/etc/ath/RADIUS_MAC.ap_bss
/etc/ath/apcfg
/lib/modules/2.6.36/atheros_wfo/asf.ko
/etc/ath/FT-EAP.ap_bss
/sbin/apdown
/etc/ath/.configData
/lib/modules/2.6.36/atheros_wfo/ath_dfs.ko
/sbin/iwpriv
/sbin/hostapd
/lib/modules/2.6.36/atheros_wfo/adf.ko
/etc/ath/PSK.ap_bss
/sbin/makeVAP
/rboot/wfo_atheros_11AC/wfo_start_dbdc.sh
/lib/libiw.so.29
/etc/ath/PSK.sta
/lib/modules/2.6.36/atheros_wfo/ath_dev.ko
/lib/modules/2.6.36/atheros_wfo/ath_rate_atheros.ko
/etc/ath/killVAP
/etc/ath/apup
/etc/ath/activateVAP
/etc/ath/wpscli
/etc/ath/WPA-FT-EAP.ap_bss
/sbin/activateVAP
/lib/modules/2.6.36/atheros_wfo/ath_hal.ko
/etc/ath/prepareACFG
/sbin/apup
/sbin/cfg
/rboot/wfo_atheros_11AC/wfo_start.sh
/etc/ath/PSK.adhoc
/etc/rc.d/rc.wlan
/etc/ath/WSC_configupdate
/etc/ath/FT-PSK.ap_bss
/sbin/iwconfig
/etc/ath/EAP.ap_bss
/etc/ath/makeVAP
/etc/ath/apdown
/etc/ath/hostapd_ctrl_iface.conf
/etc/ath/WSC.conf
/etc/config/wireless_ar9880
/lib/modules/2.6.36/atheros_wfo/ath_pktlog.ko
/sbin/wlanconfig
