# rootfs/target/common/2.6.36/etc/init.d/udhcpd
#!/bin/sh /etc/rc.common
# (C) 2008 openwrt.org

START=50

start() {
	udhcpd /etc/udhcpd.conf
}

stop() {
	pkill udhcpd
}
