# rootfs/target/origin/2.6.36/usr/lib/opkg/info/cs-tr069-be.list
/etc/config/ethport
/etc/config/tr098acclst
/etc/config/tr098notif
/etc/config/tr069id
/etc/config/tr098
/etc/init.d/ethport
/etc/config/tr069event
