# rootfs/target/origin/2.6.36/usr/lib/opkg/info/lighttpd.list
/usr/lib/lighttpd/mod_indexfile.so
/etc/lighttpd/lighttpd.conf
/etc/init.d/lighttpd
/etc/default/lighttpd
/usr/sbin/lighttpd
/usr/lib/lighttpd/mod_staticfile.so
/usr/lib/lighttpd/mod_dirlisting.so
