# rootfs/target/origin/2.6.36/usr/lib/opkg/info/NeoRaider-mrd6.list
/usr/bin/mrd6sh
/etc/mrd6.conf
/usr/sbin/mrd6
/usr/lib/mrd6/bgp.so
/etc/init.d/mrd6
