# rootfs/target/origin/2.6.36/usr/lib/opkg/info/rp-pppoe-client.list
/usr/sbin/pppoe-start
/etc/ppp/pppoe.conf
/usr/sbin/pppoe-stop
/usr/sbin/pppoe-connect
/usr/sbin/pppoe
/etc/init.d/pppoe-client
