# rootfs/target/tpgw/etc/init.d/tpgw
#!/bin/sh /etc/rc.common

START=99

TASK="node"
#OPTIONS="--stack_size=1024 --max_old_space_size=20 --max_new_space_size=4096 --max_executable_size=10 --gc_global --gc_interval=100 --trace-gc --track_gc_object_stats --trace_gc_verbose"
OPTIONS="--stack_size=1024 --max_old_space_size=20 --max_new_space_size=4096 --max_executable_size=10"
DIR="/opt/thingplus-gateway/device"
LOG_DIR="/mnt/log"
RSYNC_FILE="$DIR/update/FTM-100S/rsync.sh"

if [ -f $DIR/app.js ] ;  then
	CMD="$TASK $OPTIONS app.js"
else 
	CMD="$TASK $OPTIONS app.min.js"
fi

start() 
{
	if [ ! -e $LOG_DIR ]
	then
#		MOUNT_DEV=`cat /proc/mtd | awk '{ if ($4 ~ /log/) print $1}' | sed -e 's/://' -e 's/mtd/mtdblock/'`

#		if [ "$MOUNT_DEV" ]
#		then
			mkdir -p $LOG_DIR
#			mount -t yaffs2 /dev/$MOUNT_DEV $LOG_DIR
#		fi
	fi

    #recover rsync if needed
    $RSYNC_FILE -r

	if ! pidof $TASK | sed "s/$$\$//" | grep -q [0-9] ; then
		cd $DIR
		$CMD 2>&1 | svlogd $LOG_DIR &
	fi	
}

stop() 
{
	pkill $TASK 2> /dev/null &
}

