# rootfs/target/common/2.6.36/usr/lib/opkg/info/algapp.list
/etc/alg-rules/ipsec.rules
/etc/init.d/algapp
/etc/config/algapp
/etc/alg-rules/ftp.rules
