# rootfs/target/common/2.6.36/usr/lib/opkg/info/block-mount.list
/lib/functions/block.sh
/lib/functions/mount.sh
/etc/config/fstab
/etc/init.d/fstab
/lib/functions/fsck.sh
