# rootfs/target/common/2.6.36/usr/lib/opkg/info/busybox.list
/bin/dmesg
/usr/bin/id
/bin/login
/sbin/blockdev
/bin/rm
/usr/bin/wget
/usr/sbin/nandwrite
/sbin/lsmod
/bin/sh
/usr/bin/lsof
/sbin/vconfig
/usr/bin/test
/sbin/mkdosfs
/sbin/depmod
/bin/zcat
/bin/fgrep
/usr/sbin/killall5
/bin/ash
/sbin/halt
/sbin/reboot
/bin/true
/sbin/route
/usr/bin/md5sum
/usr/bin/free
/sbin/logread
/usr/bin/time
/usr/bin/passwd
/sbin/klogd
/usr/bin/traceroute
/bin/more
/bin/vi
/usr/bin/tftp
/bin/pwd
/sbin/modprobe
/usr/sbin/ubirmvol
/usr/bin/nc
/sbin/mke2fs
/bin/chown
/usr/bin/renice
/etc/init.d/cron
/usr/sbin/rdate
/usr/sbin/addgroup
/bin/mount
/bin/sync
/usr/sbin/adduser
/usr/bin/dirname
/sbin/rmmod
/usr/bin/hostid
/sbin/hwclock
/sbin/udhcpc
/sbin/insmod
/bin/ping6
/usr/bin/basename
/usr/sbin/delgroup
/bin/echo
/usr/bin/less
/usr/sbin/chroot
/usr/sbin/crond
/usr/bin/ftpget
/usr/bin/head
/bin/gunzip
/bin/netstat
/usr/sbin/ubimkvol
/usr/sbin/ftpd
/bin/busybox
/usr/bin/pwdx
/bin/grep
/usr/bin/clear
/usr/bin/seq
/bin/sleep
/bin/nice
/usr/bin/diff
/sbin/watchdog
/usr/bin/groups
/bin/kill
/bin/mktemp
/usr/bin/tr
/usr/bin/whois
/usr/bin/mesg
/sbin/mkfs.minix
/bin/mknod
/usr/bin/resize
/usr/sbin/nbd-client
/bin/rmdir
/usr/bin/uniq
/usr/bin/tee
/usr/bin/ftpput
/usr/sbin/powertop
/usr/bin/bzcat
/usr/bin/hexdump
/bin/date
/usr/bin/wc
/usr/bin/nslookup
/usr/bin/[
/bin/base64
/usr/bin/pgrep
/bin/umount
/sbin/swapon
/usr/bin/tail
/usr/sbin/nanddump
/bin/mpstat
/usr/bin/pkill
/usr/bin/mkfifo
/bin/chgrp
/sbin/swapoff
/bin/cp
/sbin/poweroff
/usr/bin/[[
/bin/dd
/bin/df
/bin/false
/usr/bin/sha3sum
/usr/bin/killall
/usr/bin/sort
/bin/mkdir
/sbin/ifconfig
/usr/bin/cmp
/usr/bin/reset
/sbin/start-stop-daemon
/usr/bin/pmap
/usr/sbin/dhcprelay
/bin/usleep
/bin/cat
/usr/sbin/arping
/sbin/fstrim
/sbin/mkfs.vfat
/etc/init.d/telnet
/usr/bin/telnet
/sbin/getty
/bin/touch
/sbin/mkswap
/bin/tar
/usr/bin/taskset
/usr/sbin/ubiupdatevol
/usr/bin/lsusb
/usr/sbin/ubirsvol
/bin/iostat
/usr/bin/yes
/usr/bin/top
/sbin/fdisk
/sbin/switch_root
/usr/bin/expr
/usr/bin/env
/sbin/mkfs.ext2
/sbin/sysctl
/sbin/syslogd
/bin/sed
/usr/bin/logger
/bin/chmod
/bin/setserial
/usr/bin/awk
/bin/uname
/sbin/blkid
/usr/bin/which
/usr/bin/find
/bin/ln
/usr/bin/pstree
/bin/ls
/sbin/pivot_root
/usr/bin/nohup
/usr/bin/cut
/usr/bin/du
/bin/mv
/usr/sbin/telnetd
/usr/bin/crontab
/usr/bin/printf
/usr/bin/lspci
/sbin/init
/usr/bin/tty
/bin/ping
/sbin/hdparm
/usr/bin/bunzip2
/sbin/mkfs.reiser
/usr/bin/uptime
/usr/sbin/deluser
/bin/pidof
/bin/egrep
/usr/bin/strings
/usr/sbin/udhcpd
/sbin/arp
/bin/ps
/usr/bin/xargs
/usr/bin/hd
/bin/gzip
/usr/bin/dumpleases
