# rootfs/target/common/2.6.36/usr/lib/opkg/info/davfs2.list
/usr/share/doc/davfs2/README.translators
/usr/share/doc/davfs2/INSTALL
/usr/share/doc/davfs2/THANKS
/etc/config/enter2
/usr/share/davfs2/secrets
/usr/share/doc/davfs2/NEWS
/sbin/umount.davfs
/usr/share/doc/davfs2/FAQ
/etc/init.d/davfs2
/etc/davfs2/davfs2.conf
/sbin/mount.davfs
/usr/share/doc/davfs2/BUGS
/etc/config/davfs2
/usr/sbin/mount.davfs
/lib/upgrade/keep.d/davfs2
/usr/share/doc/davfs2/AUTHORS
/usr/share/doc/davfs2/COPYING
/usr/sbin/umount.davfs
/usr/share/doc/davfs2/README
/usr/share/davfs2/davfs2.conf
/usr/share/doc/davfs2/ChangeLog
/usr/share/doc/davfs2/TODO
/etc/davfs2/secrets
