# rootfs/target/common/2.6.36/usr/lib/opkg/info/kmod-atheros-11AC-10.1.467.list
/sbin/killVAP
/lib/modules/2.6.36/asf.ko
/etc/ath/WSC_sta.conf
/etc/ath/WPA-FT-PSK.ap_bss
/lib/modules/2.6.36/ath_dfs.ko
/etc/ath/RADIUS_MAC.ap_bss
/etc/ath/apcfg
/lib/modules/2.6.36/adf.ko
/etc/ath/FT-EAP.ap_bss
/sbin/apdown
/etc/ath/.configData
/sbin/iwpriv
/lib/modules/2.6.36/ath_dev.ko
/sbin/hostapd
/lib/modules/2.6.36/ath_pktlog.ko
/etc/ath/PSK.ap_bss
/sbin/makeVAP
/lib/modules/2.6.36/umac.ko
/lib/libiw.so.29
/etc/ath/PSK.sta
/lib/modules/2.6.36/ath_hal.ko
/etc/ath/killVAP
/etc/ath/apup
/etc/ath/activateVAP
/etc/ath/wpscli
/etc/ath/WPA-FT-EAP.ap_bss
/sbin/activateVAP
/etc/ath/prepareACFG
/sbin/apup
/sbin/cfg
/etc/ath/PSK.adhoc
/etc/rc.d/rc.wlan
/lib/modules/2.6.36/ath_spectral.ko
/etc/ath/WSC_configupdate
/etc/ath/FT-PSK.ap_bss
/sbin/iwconfig
/etc/ath/EAP.ap_bss
/etc/ath/makeVAP
/etc/ath/apdown
/etc/ath/hostapd_ctrl_iface.conf
/lib/modules/2.6.36/ath_rate_atheros.ko
/etc/ath/WSC.conf
/etc/config/wireless_ar9880
/sbin/wlanconfig
