# rootfs/target/common/2.6.36/usr/lib/opkg/info/luci-server-lighttpd.list
/etc/init.d/web-server
/lib/upgrade/keep.d/lighttpd
