# rootfs/target/common/2.6.36/usr/lib/opkg/info/luci-app-minidlna.list
/etc/init.d/minidlna-gui
/usr/lib/lua/luci/controller/minidlna.lua
/etc/config/minidlna-gui
/usr/lib/lua/luci/i18n/minidlna.en.lmo
/etc/uci-defaults/luci-app-minidlna
/usr/lib/lua/luci/model/cbi/minidlna.lua
