# rootfs/target/common/2.6.36/usr/lib/opkg/info/miniupnpd.list
/etc/config/upnpd
/usr/bin/miniupnpd
/etc/hotplug.d/firewall/50-miniupnpd
/etc/init.d/miniupnpd
